# deployment/service/cherrpy/run
#!/bin/bash

exec 2>&1
exec chpst -uubuntu /home/ubuntu/git/appsforsepta/python/server.py
